# nixos/common-components/nullfs/run-nullfs.automount
[Unit]
Description=Nullfs
Before=sysinit.target

[Automount]
Where=/run/nullfs
TimeoutIdleSec=30
